# .config/systemd/user/update.service

[Unit]
Description=Update all system and AUR packages

[Service]
Type=oneshot
ExecStart=/home/ssheersh/update.sh
